# sleep_inhibitor/sleep-inhibitor.service
[Unit]
Description=Sleep Inhibitor Service

[Service]
ExecStart=/usr/bin/python3 -u /usr/bin/%p

[Install]
WantedBy=default.target
